# naive-client.service
[Unit]
Description=NaiveProxy Client (managed by naiveproxy-tui)
Documentation=https://github.com/klzgrad/naiveproxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/naiveproxy-tui/naive --config=%h/naiveproxy-tui/config.json
WorkingDirectory=%h/naiveproxy-tui
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=default.target
